# Server Tools/Mnitoring & Backups/service.template.service
[Unit]
Description="Your Service Description"

[Service]
Type=simple
ExecStart=/path/to/your/script.sh

[Install]
WantedBy=multi-user.target